# secops_rag/siem/detection_linux-threats-category.json
{
  "url": "https://docs.cloud.google.com/chronicle/docs/detection/linux-threats-category/",
  "title": "Overview of Linux Threats category",
  "content": "Home\nDocumentation\nSecurity\nGoogle Security Operations\nGuides\nStay organized with collections\nSave and categorize content based on your preferences.\nOverview of Linux Threats category\nSupported in:\nGoogle secops\nSIEM\nThis document provides an overview of the rule sets in the Linux Threats\ncategory, the required data sources, and configuration you can use to tune the\nalerts generated by these rule sets.\nRule sets in the Linux Threats category help identify threats in Linux\nenvironments using CrowdStrike Falcon, Linux Auditing System (AuditD), and Unix system logs.\nThis category includes the following rule sets:\nOS Privilege Escalation Tools\n: Detects behavior commonly seen in open\nsource Linux privilege escalation tools.\nPersistence Mechanisms:\nActivity used by adversaries to establish\nand maintain persistent access on Linux hosts.\nPrivilege Modifications\n: Activity associated with privileged\nauthentication attempts and actions, commonly used to privilege-escalate or\npersist on Linux hosts.\nMalware Signals - Suspicious LOTL Binary Activity\n: Detects suspicious\nbuilt-in tool usage (Living Off the Land) scenarios based on the observed\nactivity of Linux malware in real-world environments.\nMalware Signals - Suspicious Download Activity\n: Detects behavior seen\nin relation to malicious download activity on Linux in real-world environments.\nMalware Signals - Suspicious Execution\n: Detects signals generated from\nobserved behaviors of Linux malware detected in real-world environments\nwith a focus on Execution (TA0002) behaviors.\nMandiant Front-Line Threats -\nThis rule set contains rules derived from \nMandiant's investigation and response to active incidents across the world. \nThese rules cover commonly seen TTPs such as execution using script \ninterpreters (T1059), use of web services for command and control (T1102), \nand use of scheduled tasks to maintain persistence (T1053).\nMandiant Intel Emerging Threats -\nThis rule set contains rules derived \nfrom Mandiant Intelligence Campaigns and Significant Events, which cover \nhighly impactful geopolitical and threat activity, as assessed by Mandiant.\nThis activity may include geopolitical conflict, exploitation, phishing, \nmalvertising, ransomware, and supply chain compromises.\nSupported devices and log types\nRule sets in the Linux Threats category have been tested and are supported with\nthe following Google Security Operations supported\ndata sources\n:\nLinux Auditing System (\nAUDITD\n)\nUnix system (\nNIX_SYSTEM\n)\nCrowdStrike Falcon (\nCS_EDR\n)\nFor a list of all Google SecOps supported data sources, see\nSupported default parsers\n.\nConfigure devices to generate correct log data\nFor the rules in the Linux Threats category to function as designed, devices must\ngenerate log data in the expected format. Configure the following persistent audit rules\nfor the Linux Audit Daemon on each device where you will collect logs and sent to Google SecOps.\nFor details about how to implement persistent audit rules for the Linux Audit Daemon,\nsee documentation specific to the operating system.\n# Inserts the machine hostname into each log event\nname_format = hostname\n\n# Process creation (32 and 64-bit)\n\n-a exit,always -F arch=b32 -S execve\n-a exit,always -F arch=b64 -S execve\n\n# Persistence: Cron\n\n-w /etc/cron.allow -p wa -k cron\n-w /etc/cron.d/ -p wa -k cron\n-w /etc/cron.daily/ -p wa -k cron\n-w /etc/cron.deny -p wa -k cron\n-w /etc/cron.hourly/ -p wa -k cron\n-w /etc/cron.monthly/ -p wa -k cron\n-w /etc/cron.weekly/ -p wa -k cron\n-w /etc/crontab -p wa -k cron\n-w /var/spool/cron/ -p wa -k cron\n\n# Persistence: System Startup\n\n-w /etc/init/ -p wa -k init\n-w /etc/init.d/ -p wa -k init\n-w /etc/inittab -p wa -k init\n\n# Persistence: Systemd Units and Generators\n\n-w /etc/systemd/user -p wa -k systemd\n-w /usr/lib/systemd/user -p wa -k systemd\n-w /var/lib/systemd/linger -p wa -k systemd\n-w /root/.config/systemd/user -pa wa -k systemd\n\n-w /etc/systemd/system -p wa -k systemd\n-w /usr/lib/systemd/system -p wa -k systemd\n\n-w /run/systemd/system-generators -p wa -k systemd\n-w /etc/systemd/system-generators -p wa -k systemd\n-w /usr/local/lib/systemd/system-generators -p wa -k systemd\n-w /usr/lib/systemd/system-generators -p wa -k systemd\n\n-w /run/systemd/user-generators -pa wa -k systemd\n-w /etc/systemd/user-generators -pa wa -k systemd\n-w /usr/local/lib/systemd/user-generators -pa wa -k systemd\n-w /usr/lib/systemd/user-generators -pa wa -k systemd\n\n# Persistence: IAM\n\n-w /etc/group -p wa -k iam_etcgroup\n-w /etc/passwd -p wa -k iam_etcpasswd\n-w /etc/gshadow -k iam_etcgroup\n-w /etc/shadow -k iam_etcpasswd\n\n-w /etc/sudoers -p wa -k iam_sudoers\n-w /etc/sudoers.d/ -p wa -k iam_sudoers_d\n\n-w /usr/bin/passwd -p x -k iam_passwd\n\n-w /usr/sbin/groupadd -p x -k iam_groupmod\n-w /usr/sbin/groupmod -p x -k iam_groupmod\n-w /usr/sbin/addgroup -p x -k iam_groupmod\n-w /usr/sbin/useradd -p x -k iam_usermod\n-w /usr/sbin/userdel -p x -k iam_usermod\n-w /usr/sbin/usermod -p x -k iam_usermod\n-w /usr/sbin/adduser -p x -k iam_usermod\n\n# Privilege Escalation\n\n-w /bin/su -p x -k privesc\n-w /usr/bin/sudo -p x -k privesc\n\n# Persistence: Libraries\n\n-w /etc/ld.so.conf -p wa -k libmod\n-w /etc/ld.so.conf.d -p wa -k libmod\n-w /etc/ld.so.preload -p wa -k libmod\n\n# Persistence: PAM\n\n-w /etc/pam.d/ -p wa -k pam\n-w /etc/security/ -p wa  -k pam\n\n# Persistence: SSH\n\n-w /etc/ssh/ -p wa -k sshconfig\n-w /root/.ssh/ -p wa -k sshconfig\n\n# Persistence: Shell Configuration\n\n-w /etc/bashrc -p wa -k shellconfig\n-w /etc/csh.cshrc -p wa -k shellconfig\n-w /etc/csh.login -p wa -k shellconfig\n-w /etc/fish/ -p wa -k shellconfig\n-w /etc/profile -p wa -k shellconfig\n-w /etc/profile.d/ -p wa -k shellconfig\n-w /etc/shells -p wa -k shellconfig\n-w /etc/zsh/ -p wa -k shellconfig\n\n# Injection\n\n-a always,exit -F arch=b32 -S ptrace\n-a always,exit -F arch=b64 -S ptrace\n\n# Failed Access Attempts\n\n-a always,exit -F arch=b64 -S open -F dir=/bin -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/etc -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/home -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/sbin -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/srv -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/usr/bin -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/usr/sbin -F success=0 -k file_err\n-a always,exit -F arch=b64 -S open -F dir=/var -F success=0 -k file_err\n\n# Network connections\n\n-a always,exit -F arch=b64 -S connect -F a2=16 -F success=1 -k net_v4\n-a always,exit -F arch=b32 -S connect -F a2=16 -F success=1 -k net_v4\n-a always,exit -F arch=b64 -S connect -F a2=28 -F success=1 -k net_v6\n-a always,exit -F arch=b32 -S connect -F a2=28 -F success=1 -k net_v6\n-a always,exit -F arch=b32 -S socket -F a0=2 -k sock_v4\n-a always,exit -F arch=b64 -S socket -F a0=2 -k sock_v4\n-a always,exit -F arch=b32 -S socket -F a0=10 -k sock_v6\n-a always,exit -F arch=b64 -S socket -F a0=10 -k sock_v6\nFields needed by Linux Threats category\nThe following section describes specific data needed by rule sets in the Linux Threats category to get the greatest benefit. Make sure that your devices are configured to record the following data to device event logs.\nDataset\nUDM field (where data is stored)\nDefinition\nPrincipal Process Path\nprincipal.process.file.full_path\nLocation on disk of the current running process, if available.\nPrincipal Process Command Line\nprincipal.process.command_line\nCommand line parameters of the process, if available.\nTarget Process Path\ntarget.process.file.full_path\nLocation on disk of the target process, if available.\nTarget Process Command Line\ntarget.process.command_line\nCommand line\nNetwork DNS Query Domain\nnetwork.dns.questions.name\nDomain name of DNS queries, if available.\nTuning alerts returned by Linux Threats category\nYou can reduce the number of detections a rule or rule set generates using\nrule exclusions\n.\nIn the rule exclusion, you define the criteria of a UDM event that excludes the\nevent from being evaluated by the rule set.\nCreate one or more rule exclusions to identify criteria in a UDM event that\nexclude the event from being evaluated by this rule set or by specific rules in\nthe rule set. See\nConfigure rule exclusions\nfor information about how to do this.\nFor example, you might exclude events based on the following UDM fields:\nprincipal.hostname\ntarget.user.userid\nprincipal.process.command_line\ntarget.domain.name\ntarget.process.command_line\nNeed more help?\nGet answers from Community members and Google SecOps professionals.",
  "scraped_at": "2026-03-07T13:12:30.393583Z"
}